# home/.config/systemd/user/openrgb.service
[Unit]
Description=OpenRGB control

[Service]
Type=exec
ExecStart=/usr/bin/openrgb --server

[Install]
WantedBy=multi-user.target
